# roles/backup/files/backup.cron
0 3 * * * root /usr/local/bin/backup.sh >> /home/pi/backup.log 2>/dev/null
