# roles/backup/files/backup.cron
* 3 * * * root /usr/local/bin/backup.sh > /dev/null 2>&1
